# kubernetes/watcher/examples/kosymbiosis-exporter.service
# Systemd service for Kosymbiosis metric export
# Place in /etc/systemd/system/kosymbiosis-exporter.service

[Unit]
Description=Kosymbiosis Lambda Metric Export Service
Documentation=https://github.com/hannesmitterer/Lex-Amoris-
After=node_exporter.service

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/export-lambda.sh
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/node_exporter/textfile_collector

[Install]
WantedBy=multi-user.target
