# kubernetes/watcher/examples/node_exporter.service
# Example systemd service for node_exporter with textfile collector
# Place in /etc/systemd/system/node_exporter.service

[Unit]
Description=Prometheus Node Exporter
Documentation=https://github.com/prometheus/node_exporter
After=network.target

[Service]
Type=simple
User=node_exporter
Group=node_exporter
ExecStart=/usr/local/bin/node_exporter \
  --collector.textfile.directory=/var/lib/node_exporter/textfile_collector \
  --collector.filesystem.mount-points-exclude='^/(dev|proc|sys|var/lib/docker/.+|var/lib/kubelet/.+)($|/)' \
  --collector.netclass.ignored-devices='^(veth.*|br.*|docker.*|virbr.*|lo)$'

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
